# setup/ubuntu_docker/files/upstart.conf
description "Start re:dash Docker containers"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
    docker-compose -f /opt/redash-docker/docker-compose.yml up
end script
